# backend/crontab
0 8 * * * python /www/manage.py populate daily > populate.log
# Mandatory blank line